# src/roas2_bringup/systemd/connect-can.service
[Unit]
After=network.service time-sync.target
[Service]
Type=simple
ExecStart=/bin/sh -e /usr/sbin/connect-can
[Install]
WantedBy=multi-user.target
